# roles/statsd/templates/supervisor-statsd.conf.j2
[program:statsd]
command=/usr/local/bin/statsd /opt/statsd/config.js
process_name=%(program_name)s
autostart=true
autorestart=true
stopsignal=QUIT
