# roles/sentry/templates/supervisor-sentry.conf.j2
[program:sentry]
command=/opt/sentry/ve/bin/sentry --config=/opt/sentry/conf.py start http
environment=PATH="/opt/sentry/ve/bin"
autostart=true
autorestart=true
redirect_stderr=true
