# roles/sentry/templates/supervisor-sentry.conf.j2
[program:sentry]
command={{ venv }}/exec sentry --config={{ sentry_dir }}/conf.py start http
autostart=true
autorestart=true
redirect_stderr=true
